# install/linux/moble-tools.service
[Unit]
Description=Moble Tools (MOBI OS)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/moble-tools
Environment=HOST=0.0.0.0
Environment=PORT=4173
ExecStart=/usr/bin/node /opt/moble-tools/server.js
Restart=on-failure
RestartSec=5
User=moble
Group=moble

[Install]
WantedBy=multi-user.target
